# 10-display-dock.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", RUN+="/home/nesquiko/udev_test.sh"
